# misc/etc/cron.d/colab-wikilegis
*/30   *   *   *   * colab-admin fetch_data_wikilegis