# misc/etc/cron.d/colab-audiencias
*   *   *   *   *  colab-admin fetch_data_audiencias